# launchd/com.stevespivak.remodex.upstream.plist.template
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
  <dict>
    <key>Label</key>
    <string>com.stevespivak.remodex.upstream</string>
    <key>ProgramArguments</key>
    <array>
      <string>&lt;PYTHON_PATH&gt;</string>
      <string>-m</string>
      <string>cra.cli</string>
      <string>remodex-upstream-run</string>
      <string>--runtime-dir</string>
      <string>&lt;RUNTIME_DIR&gt;</string>
      <string>--home</string>
      <string>&lt;HOME_DIR&gt;</string>
      <string>--python-path</string>
      <string>&lt;PYTHON_PATH&gt;</string>
      <string>--node-path</string>
      <string>&lt;NODE_PATH&gt;</string>
      <string>--codex-path</string>
      <string>&lt;CODEX_PATH&gt;</string>
      <string>--remodex-bin</string>
      <string>&lt;REMODEX_BIN&gt;</string>
    </array>
    <key>WorkingDirectory</key>
    <string>&lt;REPO_ROOT&gt;</string>
    <key>EnvironmentVariables</key>
    <dict>
      <key>HOME</key>
      <string>&lt;HOME_DIR&gt;</string>
      <key>PYTHONPATH</key>
      <string>&lt;REPO_ROOT&gt;</string>
    </dict>
    <key>RunAtLoad</key>
    <true/>
    <key>KeepAlive</key>
    <true/>
    <key>ThrottleInterval</key>
    <integer>10</integer>
    <key>ProcessType</key>
    <string>Background</string>
    <key>StandardOutPath</key>
    <string>&lt;STDOUT_LOG&gt;</string>
    <key>StandardErrorPath</key>
    <string>&lt;STDERR_LOG&gt;</string>
  </dict>
</plist>
