# launchd/com.stevespivak.cra.actuator.plist.template
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
  <dict>
    <key>Label</key>
    <string>com.stevespivak.cra.actuator</string>
    <key>ProgramArguments</key>
    <array>
      <string>/usr/bin/python3</string>
      <string>-m</string>
      <string>cra.cli</string>
      <string>actuate-local</string>
      <string>--decision</string>
      <string>approve</string>
      <string>--action-id</string>
      <string>11111111-1111-4111-8111-111111111111</string>
    </array>
    <key>WorkingDirectory</key>
    <string>&lt;REPO_ROOT&gt;</string>
    <key>EnvironmentVariables</key>
    <dict>
      <key>PYTHONPATH</key>
      <string>&lt;REPO_ROOT&gt;</string>
    </dict>
    <key>RunAtLoad</key>
    <true/>
    <key>KeepAlive</key>
    <false/>
    <key>StandardOutPath</key>
    <string>&lt;REPO_ROOT&gt;/var/log/actuator.stdout.log</string>
    <key>StandardErrorPath</key>
    <string>&lt;REPO_ROOT&gt;/var/log/actuator.stderr.log</string>
  </dict>
</plist>
